# ansible/roles/deploy/templates/resolve.service.j2
[Unit]
Description=Start litecoin mainnet resolver

[Service]
WorkingDirectory=/home/ubuntu/bitnodes
ExecStart=/usr/bin/python -u resolve.py conf/resolve.conf.default

# output goes to /var/log/syslog
StandardOutput=syslog
StandardError=syslog
User=ubuntu
Group=ubuntu

SyslogIdentifier=resolve

[Install]
WantedBy=multi-user.target
Alias=resolve.service